# source_robot/ai_robot.service
[Unit]
Description=Run robot server
Wants=network-online.target
After=network-online.target

[Service]
ExecStart=/home/daan/ai_robot/venv/bin/python3 /home/daan/ai_robot/source_robot/server.py
StandardOutput=inherit
StandardInput=inherit
User=pi

[Install]
WantedBy=multi-user.target